# victron-ble/victron-ble-monitor.service
[Unit]
Description=Victron Blue Smart IP65 BLE Monitor
Documentation=https://github.com/Prinsessen/innovv-k7-autodump-
After=bluetooth.target network-online.target
Wants=bluetooth.target network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/victron-ble
ExecStart=/home/pi/victron-ble/.venv/bin/python3 /home/pi/victron-ble/victron_ble_monitor.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Give BLE stack time to initialize after boot
ExecStartPre=/bin/sleep 10

[Install]
WantedBy=multi-user.target
